# package/systemd/qsp-agent.service
[Unit]
Description=QSP Agent
Documentation=man:qsp-agent(1)
Wants=network-online.target
After=network-online.target
ConditionPathExists=/etc/qsp-agent/config.toml

[Service]
Type=simple
User=qsp-agent
Group=qsp-agent
WorkingDirectory=/var/lib/qsp-agent
RuntimeDirectory=qsp-agent
StateDirectory=qsp-agent
UMask=0027
ExecStart=/usr/bin/qsp-agent --config /etc/qsp-agent/config.toml
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
SyslogIdentifier=qsp-agent

[Install]
WantedBy=multi-user.target
